# linux/work/fs/etc/inittab
::sysinit:syslogd
::sysinit:cttyhack su -l user
::respawn:cttyhack su -l user
::ctrlaltdel:reboot
